# services/tradeflair-hourly-worker.service
[Unit]
Description=TradeFlair Hourly Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=redditbot
Group=redditbot
WorkingDirectory=/home/redditbot/Projects/TradeFlair
ExecStart=/usr/bin/pwsh -NoLogo -NoProfile -File /home/redditbot/Projects/TradeFlair/services/start-hourly-worker.ps1
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tradeflair-hourly-worker

[Install]
WantedBy=multi-user.target